# files/etc/crontabs/root
# 每天凌晨3点重启设备
0 3 * * * /sbin/reboot

# 每小时检查网络连接
0 * * * * ping -c 3 114.114.114.114 || /sbin/reboot